# roles/tilemaker/templates/tilemaker.service
[Unit]
Description=Build mbtiles

# Restart once times within 20 minutes so that an intermittent network failures are gracefully handled
StartLimitIntervalSec=1m
StartLimitBurst=3

[Service]
WorkingDirectory={{ tilemaker_work_dir }}
ExecStart={{ tilemaker_work_dir }}/build-mbtiles
RestartSec=10m
Restart=on-failure

ExecStopPost=-podman stop tilemaker
ExecStopPost=-podman rm --force tilemaker
ExecStopPost=systemd-failure-notification "tilemaker"

{% if tileserver_restart_on_tilemaker_termination %}
ExecStopPost=systemctl restart tileserver
{% endif %}

[Install]
WantedBy=multi-user.target
