# roles/photon/templates/systemd/photon-update.service
[Unit]
Description=Update photon

[Service]
ExecStart=curl localhost:{{photon_port}}/nominatim-update

[Install]
WantedBy=multi-user.target
