# roles/otp_quadlet/templates/systemd/graph-build.timer
[Unit]
Description=Build OTP graph every night

[Timer]
OnCalendar={{ graph_build_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
