# roles/otp_quadlet/templates/systemd/otp-restart.timer
[Unit]
Description=Restart otp every night

[Timer]
OnCalendar=*-*-* {{ otp_restart_time }}
Persistent=true

[Install]
WantedBy=timers.target
