# roles/tileserver/templates/systemd/tileserver-healthcheck.service
[Unit]
Description=Restart tileserver if it hangs
Requires=tileserver.service
After=tileserver.service

[Service]
RuntimeMaxSec=4m
Type=simple
ExecStart=/bin/sh -c "curl -o /dev/null --fail-with-body --retry 3 -v --max-time 30 --retry-connrefused http://localhost:{{ tileserver_port }}/styles/streets/8/135/87.png || systemctl restart tileserver"

[Install]
WantedBy=multi-user.target
