# roles/tileserver/templates/systemd/tileserver-healthcheck.timer
[Unit]
Description=Run tileserver healthcheck

[Timer]
OnCalendar={{ tileserver_healthcheck_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
