# roles/podman/templates/podman-prune.timer
[Unit]
Description=Run podman prune weekly

[Timer]
OnCalendar={{ podman_prune_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
