# roles/digitransit/templates/systemd/digitransit-restart.timer
[Unit]
Description=Restart digitransit

[Timer]
OnCalendar={{ digitransit_restart_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
